# legacy_src/worker/absen-worker.service
[Unit]
Description=Absen WhatsApp Queue Worker
After=network.target mysql.service

[Service]
Type=simple
User=www-data
# Sesuaikan path ini dengan lokasi instalasi di Ubuntu nanti
# Misalnya /var/www/html/absen/worker/worker.php
ExecStart=/usr/bin/php /var/www/html/absen/worker/worker.php
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
